# puppet/manifests/files/tilemill.conf
description "TileMill"
author      "MapBox - http://mapbox.com"

stop on shutdown
start on runlevel [2345]

respawn

script
  # See http://mapbox.com/tilemill/docs/guides/ubuntu-service/#configuring_to_listen_for_public_traffic
  su -s /bin/sh -c 'HOME=/usr/share/mapbox /usr/bin/node /usr/share/tilemill/index.js --config=/etc/tilemill/tilemill.config >> /var/log/tilemill/tilemill.log 2>&1' mapbox
end script
